# Library/LaunchAgents/rfrouter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>localhost.rfrouter</string>
        <key>ProgramArguments</key>
        <array>
            <string>sh</string>
            <string>-c</string>
            <string>
                file_path=$(readlink ~/Library/LaunchAgents/rfrouter.plist)
                router_path=$(dirname $file_path)/../../router.php
                if [ -f $router_path ]; then
                    osascript -e 'display notification "Start" with title "remote2local"'
                    # リモートマシンからのアクセスできるようにするにはlocalhostではなく0.0.0.0を使う
                    # [PHP: ビルトインウェブサーバー - Manual](http://php.net/manual/ja/features.commandline.webserver.php)
                    php -S 0.0.0.0:2224 $router_path
                else
                    osascript -e 'display notification "No such file: '$router_path'" with title "remote2local"'
                fi
            </string>
        </array>
        <key>OnDemand</key>
        <false/>
    </dict>
</plist>
